# hermes-dashboard.service
[Unit]
Description=Hermes dashboard
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m hermes.dashboard
Restart=always
RestartSec=5
WorkingDirectory=%h/.hermes
Environment="PYTHONUNBUFFERED=1"
Environment="DASHBOARD_ALLOWED_IPS=93.158.188.81,127.0.0.1"
Environment="DASHBOARD_PORT=8080"

[Install]
WantedBy=default.target
